# deploy/path-harness.service.example
[Unit]
Description=Path Harness 0.1 alpha
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/path-harness-site
Environment=HOST=0.0.0.0
Environment=PORT=8769
Environment=DB_PATH=/opt/path-harness-site/data/path_harness.sqlite
ExecStart=/usr/bin/env bash /opt/path-harness-site/scripts/run_server.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
